# systemd/home-location-endpoint.service
[Unit]
Description=Home Location Endpoint scoped Apple location interceptor
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=home-location
Group=home-location
UMask=0077
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=GSLOC_RECENT_WIFI_TTL=1800
Environment=GSLOC_RECENT_WIFI_MAX=256
Environment=GSLOC_NO_FIX_MIN_LOCATIONS=32
EnvironmentFile=-/etc/home-location-endpoint/runtime.env
ExecStart=/usr/bin/python3 /opt/home-location-endpoint/interceptor.py
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadOnlyPaths=/etc/home-location-endpoint
ReadWritePaths=/var/log/home-location-endpoint
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
SystemCallArchitectures=native
TasksMax=64
MemoryMax=256M

[Install]
WantedBy=multi-user.target
